# flask.service.j2
[Unit]
Description=Flask Web Application
After=network.target

[Service]
User=ec2-user
WorkingDirectory=/home/ec2-user
ExecStart=python3 -m flask run --host=0.0.0.0
Environment=FLASK_ENV=production

[Install]
WantedBy=multi-user.target
